# etc/systemd/system/systemd-suspend.service.d/10-no-freeze-session.conf
# /etc/systemd/system/systemd-suspend.service.d/10-no-freeze-session.conf
[Service]
Environment="SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false"
